# deploy/django-api/gunicorn/templates/best.socket
[Unit]
Description=best socket

[Socket]
ListenStream={{run_sock_file}}

[Install]
WantedBy=sockets.target